# Source/Thunder/scripts/thunder.service.in
[Unit]
Description=thunder
Wants=multi-user.target
After=multi-user.target @EXTRA_DEPENDS@

[Service]
PIDFile=@CMAKE_INSTALL_FULL_RUNSTATEDIR@/Thunder.pid
EnvironmentFile=-@CMAKE_INSTALL_FULL_SYSCONFDIR@/Thunder/Thunder.env
Environment="WAYLAND_DISPLAY=wayland-0"
Environment="XDG_RUNTIME_DIR=@CMAKE_INSTALL_FULL_RUNSTATEDIR@"
ExecStart=-@CMAKE_INSTALL_FULL_BINDIR@/Thunder -b
ExecStop=/bin/kill $MAINPID

[Install]
WantedBy=multi-user.target
